# systemd/.config/systemd/user/dotfiles-sync.service
[Unit]
Description=Dotfiles sync service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.dotfiles/scripts/dotfiles-sync.sh
Environment="DOTFILES_DIR=%h/.dotfiles"

# Don't fail the service if sync fails
SuccessExitStatus=0 1

[Install]
WantedBy=default.target
